# ipmi-facebook/service_files/xyz.openbmc_project.Chassis.Control.Ipmi.service
[Unit]
Description=Facebook Ipmi Read

[Service]
Restart=always
RestartSec=3
ExecStart=/usr/bin/fb-yv2-misc
Type=dbus
BusName=xyz.openbmc_project.State.Ipmi

[Install]
WantedBy=sysinit.target

